# lib/vendor/scripts/gentoo/solr
#!/sbin/runscript

description="Java Solr search engine."

JAVA="/usr/bin/java"

# Your symfony project
PROJECT_NAME=sfproject
PROJECT_DIR="/var/www/$PROJECT_NAME"

# Replace with the solr installation path
SOLR_HOME="$PROJECT_DIR/plugins/tjSolrDoctrineBehaviorPlugin/lib/vendor/solr"

# The data directory
SOLR_DATA="$PROJECT_DIR/data/solr"

# The log directory
SOLR_LOG="$PROJECT_DIR/log/solr"

SOLR_OPTIONS="-Dsolr.data.dir=$SOLR_DATA"
SOLR_OPTIONS="$PARAMETERS -Djava.util.logging.config.file=$SOLR_HOME/logging.properties"

SOLR_PID="/var/run/solr.pid"

start() {
	ebegin "Starting solr"
	start-stop-daemon --start --background --quiet --user solr --chdir $SOLR_HOME \
		--make-pidfile --pidfile $SOLR_PID --exec $JAVA -- $SOLR_OPTIONS -jar start.jar
	eend $?
}

stop() {
	ebegin "Stopping solr"
	start-stop-daemon --stop --quiet --pidfile $SOLR_PID
	eend $?
}
